# install/integrations/os/com.crabcc.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crabcc.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.cargo/bin/crabcc</string>
        <string>--mcp-http</string>
        <string>127.0.0.1:8091</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/Crabcc/mcp.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/Crabcc/mcp.stderr.log</string>
</dict>
</plist>
