# install/integrations/os/crabcc-mcp.service
[Unit]
Description=crabcc MCP server (stdio bridge via socket activation placeholder)
Documentation=https://github.com/peterlodri-sec/crabcc
After=network-online.target

[Service]
Type=simple
Environment=PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin
ExecStart=%h/.cargo/bin/crabcc --mcp-http 127.0.0.1:8091
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
