# installer/Crabcc.app/Contents/Resources/com.crabcc.manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crabcc.manager</string>

    <!-- Bulletproof orchestrator: long-running heartbeat + reaper.
         KeepAlive=true so any crash respawns; ProcessType=Background
         keeps it out of the way of user-interactive work. -->
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.cargo/bin/crabcc</string>
        <string>manager</string>
        <string>daemon</string>
        <string>--interval</string>
        <string>30</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>Nice</key>
    <integer>5</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOME__/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/Crabcc/manager.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/Crabcc/manager.stderr.log</string>
</dict>
</plist>
