# installer/Crabcc.app/Contents/Resources/com.crabcc.agentd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crabcc.agentd</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__APP_PATH__/Contents/Resources/scripts/crabcc-agentd.sh</string>
    </array>

    <!-- Survives sleep/wake/restart — launchd auto-restarts on crash and
         relaunches at login. ProcessType=Background lowers QoS so the
         agent yields to user-interactive work. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Throttle restart loops so a broken binary can't pin a CPU core. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Resource ceilings. -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>
    <key>Nice</key>
    <integer>5</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CRABCC_BIN</key>
        <string>__HOME__/.cargo/bin/crabcc</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/Crabcc/agentd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/Crabcc/agentd.stderr.log</string>
</dict>
</plist>
